# contrib/linux-packaging/archlinux/neucoind.service
[Unit]
Description=EmollientCoin daemon
After=network.target

[Service]
ExecStart=EmollientCoind -logtimestamps

[Install]
WantedBy=default.target
